# .config/crontab/laravel-cronjob
SHELL=/bin/bash
COLUMNS=320
* * * * * root source /var/app/current/export-env.sh && php /var/app/current/artisan schedule:run >> /dev/null 2>&1
